# launchd/com.local.chatgpt-backup-automation.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.chatgpt-backup-automation</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/one/chatgpt-backup-automation/scripts/run-scheduled.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/one/chatgpt-backup-automation</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/one/chatgpt-backup-automation/.local/logs/launchagent.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/one/chatgpt-backup-automation/.local/logs/launchagent.err.log</string>
</dict>
</plist>
